# systemd/kfxbrain-api.service
[Unit]
Description=Kurage FX Brain API
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/kojima/work/kfxbrain
EnvironmentFile=/home/kojima/work/kfxbrain/.env
ExecStart=/home/kojima/work/kfxbrain/.venv/bin/uvicorn kfxbrain.api:app --host 0.0.0.0 --port 18326
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
